# etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

. /lib/functions/gl_util.sh
/etc/firewall.nodogsplash # Set firewall rules after gl_util.sh
remount_ubifs

exit 0
