# packages/base-files/files/etc/init.d/99-first-login
#!/bin/sh /etc/rc.common

START=99
STOP=15

start() {
    if [ ! -f /etc/first_login_done ]; then
        # Your custom commands here
        
        # Mark as done
        touch /etc/first_login_done
    fi
}

stop() {
    # Optional: Add any cleanup steps here
    :
}